# root/etc/s6-overlay/s6-rc.d/apache2-frontend/run
#!/command/execlineb -P
with-contenv
s6-setuidgid root
httpd -f /etc/apache2/httpd-frontend.conf -DFOREGROUND
